# config/quadlet/certbot.timer
[Timer]
Unit = certbot.service
# Twice a day
OnCalendar = *-*-* 0,12:00:00
